# images/ubuntu/config/supervisord.conf
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL for a unix socket

[unix_http_server]
file = /tmp/supervisor.sock   ; (the path to the socket file)
chmod = 0760
chown = www-data:www-data
;username = cec
;password = 441c31e22e55930ed68c08383535

[supervisord]
user = root
logfile = /tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes = 50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups = 10           ; (num of main logfile rotation backups;default 10)
loglevel = info                ; (log level;default info; others: debug,warn,trace)
pidfile = /tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon = false               ; (start in foreground if true;default false)
minfds = 1024                  ; (min. avail startup file descriptors;default 1024)
minprocs = 200                 ; (min. avail process descriptors;default 200)
umask = 022                    ; (process file creation umask;default 022)
identifier = supervisor        ; (supervisord identifier;default is 'supervisor')
directory = /tmp               ; (default is not to cd during start)
nocleanup = true               ; (don't clean up tempfiles at start;default false)
childlogdir = /tmp             ; ('AUTO' child log dir, default $TEMP)
strip_ansi = false             ; (strip ansi escape codes in logs; def. false)


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:php-fpm]
command = /usr/sbin/php-fpm8.3 --nodaemonize --fpm-config /etc/php/8.3/fpm/php-fpm.conf
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:crond]
command = /usr/sbin/cron -f
autostart = true
autorestart = true
user = root  ; Ensure crond runs as root
stdout_logfile = /var/log/crond_stdout.log
stderr_logfile = /var/log/crond_stderr.log
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
pidfile = /var/run/cron.pid   ; Specify an alternative PID file path


[program:beanstalkd]
command = /usr/bin/beanstalkd
autostart = true
autorestart = true
user = www-data
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:queue-worker]
process_name = %(program_name)s_%(process_num)02d
command = php /var/www/artisan queue:work --sleep=3 --tries=3 --timeout=2000
autostart = true
autorestart = true
user = www-data
numprocs = 16
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0

# ENTERPRISE WEBHOOK WORKERS - SaaS Scale
# High Priority Webhook Workers - Critical tenant operations
[program:webhook-high-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-high --sleep=1 --tries=3 --timeout=300 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Beanstalk Webhook Workers - Fast processing for medium-high priority
[program:webhook-beanstalk-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-beanstalk --sleep=2 --tries=3 --timeout=300 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Redis Webhook Workers - Immediate processing
[program:webhook-redis-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-redis --sleep=2 --tries=3 --timeout=300 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Database Webhook Workers - Reliable fallback
[program:webhook-db-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-db --sleep=3 --tries=3 --timeout=300 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Low Priority Webhook Workers - Bulk operations, reports
[program:webhook-low-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-low --sleep=5 --tries=5 --timeout=600 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Failed Webhook Retry Workers - Handle failed webhook retries
[program:webhook-failed-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work webhook-failed --sleep=10 --tries=5 --timeout=600 --memory=512
directory=/var/www
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:socket-worker]
process_name = %(program_name)s_%(process_num)02d
command = php /var/www/artisan websockets:serve
autostart = true
autorestart = true
user = www-data
numprocs = 1
redirect_stderr = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0

# Worker Groups for easier management
[group:core-services]
programs=nginx,php-fpm,crond,beanstalkd
priority=999

[group:general-workers]
programs=queue-worker,socket-worker
priority=998

[group:webhook-workers]
programs=webhook-high-worker,webhook-beanstalk-worker,webhook-redis-worker,webhook-db-worker,webhook-low-worker,webhook-failed-worker
priority=997

[inet_http_server]
port = *:9001
username = admin
password = admin