# images/gateway/setasign.7.4/config/crontab
# Tenant
* * * * * cd /var/www && php artisan tenancy:run schedule:run >> /dev/null 2>&1

# System
* * * * * cd /var/www && php artisan schedule:run >> /dev/null 2>&1
